# timers/msi_backup.timer
[Unit]
Description=Run MSI Backup weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target