# timers/msi_backup.service
[Unit]
Description=Run MSI Backup

[Service]
ExecStart=/usr/local/bin/msi_backup.sh