# timers/sysupdate.service
[Unit]
Description=Run system update

[Service]
ExecStart=/usr/local/bin/system-update.sh
